# systemd/dockerpost-learnspam.timer
[Unit]
Description=Timer for training SpamAssassin's Bayesian filter

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
